# target/linux/x86/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
::respawn:/sbin/getty 115200 tty1
ttyS1::askfirst:/bin/ash --login
